# scripts/ec2/systemd/heptad-relayer-sepolia@.service
[Unit]
Description=Heptad Sepolia Relayer (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ssm-user
WorkingDirectory=/home/ssm-user/heptad
Environment=HOME=/home/ssm-user
Environment=NODE_VERSION=24
EnvironmentFile=/home/ssm-user/heptad-env/relayer-common.env
EnvironmentFile=-/home/ssm-user/heptad-env/sepolia-%i.env
ExecStart=/bin/bash -lc '. "$HOME/.nvm/nvm.sh" && nvm use "$NODE_VERSION" >/dev/null && npm run --workspace @heptad/relayer run:sepolia'
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
