# scripts/ec2/systemd/heptad-relayer-aggregate-mint.service
[Unit]
Description=Heptad Relayer Mint Aggregator (DB -> API snapshot)
After=network-online.target heptad-relayer-api.service
Wants=network-online.target

[Service]
Type=oneshot
User=ssm-user
WorkingDirectory=/home/ssm-user/heptad
Environment=HOME=/home/ssm-user
Environment=NODE_VERSION=24
Environment=RELAYER_ATTESTATION_SOURCE=db
Environment=RELAYER_DISABLE_FILE_OUTPUT=1
EnvironmentFile=/home/ssm-user/heptad-env/relayer-common.env
EnvironmentFile=-/home/ssm-user/heptad-env/relayer-api.env
EnvironmentFile=-/home/ssm-user/heptad-env/aggregator.env
ExecStart=/bin/bash -lc '. "$HOME/.nvm/nvm.sh" && nvm use "$NODE_VERSION" >/dev/null && npm run --workspace @heptad/relayer aggregate:sepolia'

[Install]
WantedBy=multi-user.target
