# scripts/ec2/systemd/opbridge-relayer-api.service
[Unit]
Description=OP_BRIDGE Relayer API (status + candidate registry)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ssm-user
WorkingDirectory=/home/ssm-user/opbridge
Environment=HOME=/home/ssm-user
Environment=NODE_VERSION=24
Environment=RELAYER_API_PORT=8787
Environment=RELAYER_API_DB_PATH=/home/ssm-user/opbridge/services/api/.data/relayer-api.sqlite
EnvironmentFile=-/home/ssm-user/opbridge-env/relayer-api.env
ExecStartPre=/bin/bash -lc '. "$HOME/.nvm/nvm.sh" && nvm use "$NODE_VERSION" >/dev/null && npm run --workspace @opbridge/api import:relayer-snapshots'
ExecStart=/bin/bash -lc '. "$HOME/.nvm/nvm.sh" && nvm use "$NODE_VERSION" >/dev/null && npm run --workspace @opbridge/api start'
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
