# scripts/ec2/systemd/opbridge-relayer-ethereum@.service
[Unit]
Description=OP_BRIDGE Ethereum Relayer (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ssm-user
WorkingDirectory=/home/ssm-user/opbridge
Environment=HOME=/home/ssm-user
Environment=NODE_VERSION=24
EnvironmentFile=/home/ssm-user/opbridge-env/relayer-common.env
EnvironmentFile=-/home/ssm-user/opbridge-env/ethereum-%i.env
ExecStart=/bin/bash -lc '. "$HOME/.nvm/nvm.sh" && nvm use "$NODE_VERSION" >/dev/null && npm run --workspace @opbridge/relayer run:ethereum'
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
